# 实时同步/sync.ini
[program: sync_cofig]
command=/bin/sh /data/scripts/sync/inotify_monitor.sh
priority=10
autostart=true
autorestart=true
startsecs=10
startretries=3

